# etc/init.d/touchscreen
#!/sbin/openrc-run

name="touchscreen"
command="/system/bin/run-rm-wrapper"
